# debian-control/postinst.in
#!/usr/bin/env bash

set +e

useradd -r nanocurrency
groupadd -r nanocurrency

mkdir -p /var/nanocurrency/Vban
mkdir -p /var/nanocurrency/VbanBeta
mkdir -p /var/nanocurrency/VbanTest
chown -R nanocurrency: /var/nanocurrency

systemctl stop ${VBAN_SERVICE} 2>/dev/null
systemctl disable ${VBAN_SERVICE} 2>/dev/null
cp ${CMAKE_INSTALL_PREFIX}/extras/systemd/${VBAN_SERVICE} /etc/systemd/system/${VBAN_SERVICE}
systemctl daemon-reload
systemctl enable ${VBAN_SERVICE}
systemctl start ${VBAN_SERVICE}
sleep 1
systemctl --no-pager status ${VBAN_SERVICE}

set -e
